# systemd/seem-agent.service
[Unit]
Description=SEEM Sovereign Cognitive Microservice
After=network.target

[Service]
# IMPORTANT: Update the path below to match YOUR actual installation directory
ExecStart=/usr/bin/python3 /home/SEEM/SEEM-2.0-Self-Evolving-Emergent-Mind/seem.py daemon
WorkingDirectory=/home/brian/SEEM-2.0-Self-Evolving-Emergent-Mind

# Run as your normal user (not root)
User=
Group=

# Restart policy
Restart=always
RestartSec=5

# Environment variables
Environment="API_KEY=your-secure-vsa-key-123"
Environment="PYTHONUNBUFFERED=1"

# Optional resource limits (uncomment if needed on low-power hardware)
# MemoryMax=512M
# CPUQuota=60%

[Install]
WantedBy=multi-user.target